# scripts/kiosk-client/portable-kiosk.service.template
[Unit]
Description=Portable Kiosk Launcher
After=multi-user.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@KIOSK_USER@
Group=@KIOSK_USER@
WorkingDirectory=@INSTALL_DIR@
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/@KIOSK_USER@/.Xauthority
ExecStartPre=/bin/sleep 10
ExecStart=@INSTALL_DIR@/venv/bin/python @INSTALL_DIR@/portable_launcher.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
